# supervisor.conf
[inet_http_server]
port=127.0.0.1:9090

[program:vncserver]
command=tigervncserver -fg -geometry 1600x900 :1 -- startlxde
stdout_logfile=/var/log/vncserver.log
redirect_stderr=true
autostart=true

[program:noVNC]
command=/noVNC/utils/launch.sh --vnc localhost:5901
stdout_logfile=/var/log/novnc.log
redirect_stderr=true

[program:Minio]
command=/opt/minio server --console-address :37173 /root/Desktop/data/
environment=MINIO_ROOT_USER="admin",MINIO_ROOT_PASSWORD="password",MINIO_ETCD_PATH_PREFIX="s3/",MINIO_ACCESS_KEY="key",MINIO_SECRET_KEY="secret"
stdout_logfile=/var/log/minio.log
redirect_stderr=true
